# cf-ddns.timer
[Unit]
Description=Run on boot

[Timer]
OnBootSec=0sec

[Install]
WantedBy=timers.target
